# deploy/phoromatic-upstart/phoromatic-server.conf
# description     "Phoronix Test Suite Phoromatic Server"
# author          "Phoronix Media"

start on (local-filesystems and net-device-up IFACE!=lo)
stop on stopping network
stop on starting shutdown

respawn

exec /usr/bin/phoronix-test-suite start-phoromatic-server
